# gnome/dbus/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/dbus-1/
-rw-r--r--	root/root	etc/dbus-1/session.conf
-rw-r--r--	root/root	etc/dbus-1/system.conf
drwxr-xr-x	root/root	etc/dbus-1/system.d/
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/dbus
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/dbus-cleanup-sockets
-rwxr-xr-x	root/root	usr/bin/dbus-daemon-1
-rwxr-xr-x	root/root	usr/bin/dbus-glib-tool
-rwxr-xr-x	root/root	usr/bin/dbus-launch
-rwxr-xr-x	root/root	usr/bin/dbus-monitor
-rwxr-xr-x	root/root	usr/bin/dbus-send
-rwxr-xr-x	root/root	usr/bin/dbus-viewer
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/dbus-1.0/
drwxr-xr-x	root/root	usr/include/dbus-1.0/dbus/
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/connection.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-address.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-bus.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-connection.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-errors.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-glib-error-enum.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-glib-lowlevel.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-glib.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-macros.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-memory.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-message.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-pending-call.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-protocol.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-qt.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-server.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-shared.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-threads.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus-types.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/dbus.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/message.h
-rw-r--r--	root/root	usr/include/dbus-1.0/dbus/server.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/dbus-1.0/
drwxr-xr-x	root/root	usr/lib/dbus-1.0/include/
drwxr-xr-x	root/root	usr/lib/dbus-1.0/include/dbus/
-rw-r--r--	root/root	usr/lib/dbus-1.0/include/dbus/dbus-arch-deps.h
drwxr-xr-x	root/root	usr/lib/dbus-1.0/services/
-rw-r--r--	root/root	usr/lib/libdbus-1.a
-rwxr-xr-x	root/root	usr/lib/libdbus-1.la
lrwxrwxrwx	root/root	usr/lib/libdbus-1.so -> libdbus-1.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libdbus-1.so.0 -> libdbus-1.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libdbus-1.so.0.0.0
-rw-r--r--	root/root	usr/lib/libdbus-glib-1.a
-rwxr-xr-x	root/root	usr/lib/libdbus-glib-1.la
lrwxrwxrwx	root/root	usr/lib/libdbus-glib-1.so -> libdbus-glib-1.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libdbus-glib-1.so.0 -> libdbus-glib-1.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libdbus-glib-1.so.0.0.0
-rw-r--r--	root/root	usr/lib/libdbus-qt-1.a
-rwxr-xr-x	root/root	usr/lib/libdbus-qt-1.la
lrwxrwxrwx	root/root	usr/lib/libdbus-qt-1.so -> libdbus-qt-1.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libdbus-qt-1.so.1 -> libdbus-qt-1.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libdbus-qt-1.so.1.0.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/dbus-1.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/dbus-glib-1.pc
drwxr-xr-x	root/root	usr/lib/python2.3/
drwxr-xr-x	root/root	usr/lib/python2.3/site-packages/
-rw-r--r--	root/root	usr/lib/python2.3/site-packages/dbus.py
-rw-r--r--	root/root	usr/lib/python2.3/site-packages/dbus.pyc
-rw-r--r--	root/root	usr/lib/python2.3/site-packages/dbus.pyo
-rw-r--r--	root/root	usr/lib/python2.3/site-packages/dbus_bindings.a
-rwxr-xr-x	root/root	usr/lib/python2.3/site-packages/dbus_bindings.la
-rwxr-xr-x	root/root	usr/lib/python2.3/site-packages/dbus_bindings.so
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/dbus-cleanup-sockets.1.gz
-rw-r--r--	root/root	usr/man/man1/dbus-daemon-1.1.gz
-rw-r--r--	root/root	usr/man/man1/dbus-launch.1.gz
-rw-r--r--	root/root	usr/man/man1/dbus-monitor.1.gz
-rw-r--r--	root/root	usr/man/man1/dbus-send.1.gz
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxr-xr-x	messagebus/root	var/lib/dbus/
drwxr-xr-x	root/root	var/run/
drwxr-xr-x	messagebus/root	var/run/dbus/
